# conf/supervise-gazetteer.conf
[program:gazetter]
command=/bin/bash -c "/lucene-geo-gazetteer/src/main/bin/lucene-geo-gazetteer -server"
enabled=true
autorestart=false

redirect_stderr = True
redirect_stdout = True
